# exit/etc/profile.d/conda.csh
setenv CONDA_EXE "/Users/vanessadungl/Desktop/AI-Energy-Dashboard/exit/bin/conda";
setenv _CONDA_ROOT "/Users/vanessadungl/Desktop/AI-Energy-Dashboard/exit";
setenv _CONDA_EXE "/Users/vanessadungl/Desktop/AI-Energy-Dashboard/exit/bin/conda";
setenv CONDA_PYTHON_EXE "/Users/vanessadungl/Desktop/AI-Energy-Dashboard/exit/bin/python";
source "/Users/vanessadungl/Desktop/AI-Energy-Dashboard/exit/lib/python3.12/site-packages/conda/shell/etc/profile.d/conda.csh";
